# dotfiles/config/autostart/xrandr-aml.desktop
[Desktop Entry]
Name=Fix AML Monitors
Comment=Runs xrandr setup for AML monitors on login
Exec=xrandr-aml
